# noolite-mqtt.service
[Unit]
Description=Noolite control via mqtt

[Service]
ExecStart=/usr/bin/noolite_mqtt &

[Install]
WantedBy=multi-user.target

